# systemd/codex-usagebar.service
[Unit]
Description=Update Codex usage status for i3status

[Service]
Type=oneshot
ExecStart=%h/.local/bin/codex-usagebar
