# systemd/codex-usagebar.timer
[Unit]
Description=Refresh Codex usage status for i3status

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=10s
Unit=codex-usagebar.service

[Install]
WantedBy=timers.target
