# systemd/app-dev.lizardbyte.app.Sunshine.service.d/ipad-display.conf
[Service]
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
ExecStartPre=
ExecStartPre=/bin/sleep 2
ExecStartPre=%h/bin/ipad-display-apply
